# steamlink.service
[Unit]
Description=Steamlink Service

Conflicts=display-manager.service

OnFailure=display-manager.service

[Service]
User=pi
Group=pi
ExecStart=/usr/bin/steamlink

Environment="XDG_RUNTIME_DIR=/run/user/1000" # TODO: check, should probably be dynamic

ExecStopPost=sudo service lightdm start

[Install]
WantedBy=default.target
