# systemd/bluetooth-devices-connector.service
# NOTE: in order to inspect service output
# $ journalctl -e -b -u bluetooth-devices-connector.service

[Unit]
Description=Execute script.bluetooth-devices-connector kodi addon

[Service]
Type=oneshot

ExecStart=/usr/bin/python /home/pi/.kodi/addons/script.bluetooth-devices-connector/cron.py

[Install]
WantedBy=multi-user.target
