# systemd/steamlink.service
[Unit]
Description=Steamlink Service

Conflicts=display-manager.service

OnFailure=display-manager.service

[Service]
User=pi

Environment="XDG_RUNTIME_DIR=/run/user/1000"

ExecStartPre=sudo service lightdm stop

ExecStart=sudo -u pi -E /usr/bin/steamlink

ExecStopPost=sudo service lightdm start

[Install]
WantedBy=default.target
