# systemd/user/bluetooth-devices-connector.service
[Unit]
Description=Bluetooth devices auto-connector

[Service]
Type=simple

ExecStart=/bin/bash -c "/bin/sh /home/pi/bin/connect_bluetooth_devices.sh \"$(< /home/pi/.kodi/userdata/addon_data/script.bluetooth-devices-connector/devices)\""

[Install]
WantedBy=default.target
